# airootfs/etc/autostart/kscreen.desktop
[Desktop Entry]
Type=Application
Exec=kscreen_backend_launcher
Name=KScreen
